# overlay/system/etc/init/dt2w.conf
description "Enable DT2W on boot."

start on started lightdm

task

script
        echo 1 > /proc/touchpanel/double_tap_enable
end script
