# overlay/system/etc/init/repowerd_fix.conf
description "Fix autobrightness"

start on started lightdm

task

script
        service repowerd restart
end script

