# overlay/system/etc/init/set_props.conf
description "Set required additional props"

start on starting lightdm

task

script
        setprop ubuntu.unity8.interactive_blur true
end script
